# ultimate-agent.service
[Unit]
Description=Ultimate Agent - Unified Python Edition
After=network.target ollama.service

[Service]
Type=simple
User=zeds
WorkingDirectory=/home/zeds/Desktop/ultimate-agent/python-agent
Environment="PATH=/home/zeds/Desktop/ultimate-agent/python-agent/venv/bin"
# PYTHONPATH is needed for app imports to work correctly
Environment="PYTHONPATH=/home/zeds/Desktop/ultimate-agent/python-agent"
ExecStart=/home/zeds/Desktop/ultimate-agent/python-agent/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
